# test.rules
ACTION=="add" ,KERNEL=="sdb*", RUN+="/home/user/simran/my_script.sh"
